# supervisor.apps.ini
[program:example]
directory = /flask
command = venv/bin/gunicorn
autostart = true
autorestart = false
startretries = 0
